# src/MCPServer/debian/postinst
#!/bin/sh

set -e

#Commented out until further tests are completed
#. /usr/share/debconf/confmodule
#. /usr/share/dbconfig-common/dpkg/postinst 

#dbc_dbname=MCP
#dbc_dbuser=demo

#dbc_go archivematica-mcp-server $@

sudo adduser --uid 333 --group --system --home /var/lib/archivematica/ archivematica
chown -R archivematica:archivematica "/var/archivematica/"
chmod -R g+s "/var/archivematica/"
chmod -R 775 "/var/archivematica/"
echo "archivematica ALL=NOPASSWD:/bin/mv,/bin/chown,/bin/chmod,/usr/bin/gs,/usr/lib/transcoder/transcoderScripts/DocumentConverter.py" >> /etc/sudoers

mkdir /var/log/archivematica
chown archivematica:archivematica /var/log/archivematica

if [ -x "/etc/init.d/archivematica-mcp-serverd" ]; then
        update-rc.d archivematica-mcp-serverd defaults >/dev/null
        #invoke-rc.d archivematica-mcp-serverd start || exit $?
fi
